# docker/web/s6-rc.d/pnpm-watcher/run
#!/command/with-contenv bash
# shellcheck shell=bash
# Watch pnpm-lock.yaml for changes and sync node_modules automatically.
# Uses polling (3s interval) — inotify events don't propagate from the host
# through bind mounts on this setup (see WATCHPACK_POLLING in docker-compose).

LOCKFILE="/app/pnpm-lock.yaml"
SENTINEL="/app/node_modules/.pnpm-install-stamp"
POLL_INTERVAL=3

log() { echo "[pnpm-watcher] $*" >&2; }

sync_packages() {
    log "pnpm-lock.yaml changed — running pnpm install"
    if cd /app && CI=true s6-setuidgid node pnpm install --frozen-lockfile >&2; then
        touch "$SENTINEL"
        log "done — stopping pnpm-dev"
        s6-svc -d /run/service/pnpm-dev
        # Wait until s6 confirms the service's run script has exited.
        # s6-svwait -D returns as soon as pnpm-dev is in "down" state.
        /command/s6-svwait -D -t 15000 /run/service/pnpm-dev 2>/dev/null || true
        # Port guard: poll until :49010 is free. The finish script kills the orphan
        # concurrently; this loop waits for it regardless of how long that takes.
        i=0
        while [ "$i" -lt 15 ] && ss -tlnp 2>/dev/null | grep -q ':49010'; do
            sleep 1
            i=$((i + 1))
        done
        log "starting pnpm-dev"
        s6-svc -u /run/service/pnpm-dev
    else
        log "pnpm install failed"
    fi
}

if [ ! -f "$LOCKFILE" ]; then
    log "pnpm-lock.yaml not found — exiting"
    exit 0
fi

log "polling $LOCKFILE every ${POLL_INTERVAL}s"

last_mtime=$(stat -c %Y "$LOCKFILE" 2>/dev/null)

while true; do
    sleep "$POLL_INTERVAL"
    mtime=$(stat -c %Y "$LOCKFILE" 2>/dev/null)
    if [ -n "$mtime" ] && [ "$mtime" != "$last_mtime" ]; then
        last_mtime="$mtime"
        sync_packages
    fi
done
