# docker/web/s6-rc.d/pnpm-dev/finish
#!/bin/sh
# Called by s6 after pnpm-dev exits (before restart).
# pnpm run dev spawns next dev which spawns turbopack workers — these become
# orphan processes when pnpm exits and keep holding :49010, causing EADDRINUSE
# on the next start. Kill them here before s6 restarts the service.

# Primary: kill by process name (requires procps package)
pkill -9 -f "next" 2>/dev/null || true

# Fallback: kill whatever process is still holding :49010 directly.
# Catches turbopack workers and other node processes that don't match "next".
port_pid=$(ss -tlnp 2>/dev/null | sed -n 's/.*:49010.*pid=\([0-9]*\).*/\1/p')
[ -n "$port_pid" ] && kill -9 "$port_pid" 2>/dev/null || true

# Poll until the port is released (up to 10s) — more reliable than a fixed sleep.
i=0
while [ "$i" -lt 10 ] && ss -tlnp 2>/dev/null | grep -q ':49010'; do
    sleep 1
    i=$((i + 1))
done
