# docker/web/s6-rc.d/pnpm-dev/run
#!/command/with-contenv bash
# shellcheck shell=bash
# Next.js dev server. Source is bind-mounted at /app at runtime.
# next dev binds to all interfaces when HOST=0.0.0.0 (set in docker-compose env).
cd /app
exec s6-setuidgid node pnpm run dev
