# docker/s6/services.d/batch-worker/run
#!/command/with-contenv bash
exec s6-setuidgid axon /usr/local/bin/cortex batch worker
